# supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/var/log/supervisord/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
user=appuser
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
autorestart=true
priority=100
startsecs=5
stopwaitsecs=10

[program:backend]
command=python3 /app/backend/run.py
directory=/app/backend
stdout_logfile=/var/log/backend/app.log
stderr_logfile=/var/log/backend/error.log
user=appuser
autorestart=true
priority=200
startsecs=10
stopwaitsecs=10
environment=PYTHONUNBUFFERED="1",FLASK_ENV="production",PROTOCOL_BUFFERS_PYTHON_IMPLEMENTATION=python,DATABASE_URL=sqlite+libsql:////app/data/local.db

[group:coredns-admin]
programs=nginx,backend
priority=999